# systemd/gitlab-backup.timer
[Unit]
Description=Runs gitlab-backup once a day

[Timer]
OnCalendar=daily
RandomizeDelaySec=7200